# deploy/tradeclaw.service
[Unit]
Description=TradeClaw Trading Bot
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/root/TradeClaw
EnvironmentFile=/root/TradeClaw/.env
ExecStart=/usr/bin/npx tsx src/main.ts
Restart=on-failure
RestartSec=10
# Logging
StandardOutput=append:/root/TradeClaw/app.log
StandardError=append:/root/TradeClaw/app.log

# Graceful shutdown
KillSignal=SIGTERM
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
